# src/main/assembly/dist/bin/easy-pid-generator
#!/usr/bin/env bash
# chkconfig: 2345 91 59

### BEGIN INIT INFO
# Provides:          easy-pid-generator
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Short-Description: Starts the easy-pid-generator service
# Description:       This file is used to start the daemon
#                    and should be placed in /etc/init.d
### END INIT INFO

NAME="easy-pid-generator"
EXEC="/usr/bin/jsvc"
APPHOME="/opt/easy-pid-generator"
JAVA_HOME="/usr/lib/jvm/jre"
CLASSPATH="$APPHOME/bin/$NAME.jar:`echo $APPHOME/lib/*.jar | sed 's/ /:/g'`"
CLASS="nl.knaw.dans.easy.pid.ServiceStarter"
ARGS=""
USER="easy_pid_generator"
PID="/var/run/$NAME.pid"
OUTFILE="/var/log/$NAME/$NAME.out"
ERRFILE="/var/log/$NAME/$NAME.err"
WAIT_TIME=20

jsvc_exec()
{
    cd $APPHOME
    $EXEC -home $JAVA_HOME -cp $CLASSPATH -user $USER -outfile $OUTFILE -errfile $ERRFILE -pidfile $PID -wait $WAIT_TIME \
          -Dapp.home=$APPHOME -Dconfig.file=$APPHOME/cfg/application.conf \
          -Dlogback.configurationFile=$APPHOME/cfg/logback.xml $1 $CLASS $ARGS
}

start_jsvc_exec()
{
    echo "Starting $NAME ..."
    jsvc_exec
    if [[ $? == 0 ]]; then
        echo "$NAME has started."
    else
        echo "$NAME did not start successfully (exit code: $?)"
    fi
}

stop_jsvc_exec()
{
    echo "Stopping $NAME ..."
    jsvc_exec "-stop"
    if [[ $? == 0 ]]; then
        echo "$NAME has stopped."
    else
        echo "$NAME did not stop successfully (exit code: $?)"
    fi
}

case "$1" in
    start)
        start_jsvc_exec
    ;;
    stop)
        stop_jsvc_exec
    ;;
    restart)
        if [ -f "$PID" ]; then
            echo "Restarting $NAME ..."
            jsvc_exec "-stop"
            jsvc_exec
            echo "$NAME has restarted."
        else
            echo "$NAME is not running, no action taken"
            exit 1
        fi
    ;;
    status)
        if [ -f "$PID" ]
        then
            echo "$NAME (pid `cat $PID`) is running"
        else
            echo "$NAME is stopped"
        fi
    ;;
    *)
    echo "Usage: sudo service $NAME {start|stop|restart|status}" >&2
    exit 3
esac
